# codex-ollama-proxy.service
[Unit]
Description=Codex ↔ Ollama Proxy — OpenAI API to Ollama bridge
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
ExecStart=/usr/bin/node /home/grey/codex-ollama-proxy/proxy.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codex-ollama-proxy
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
