# ollama-proxy-user.service
[Unit]
Description=Ollama Proxy for Codex - Transforms developer roles to system roles
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/codex-brain/codex-ollama-proxy
ExecStart=/usr/bin/node %h/codex-brain/codex-ollama-proxy/proxy.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama-proxy

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
